# web-ui/systemd/gradient-api.service
[Unit]
Description=GradientOS HTTP API
After=network.target

StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple

User=pi
Group=pi

WorkingDirectory=/home/pi/src/GradientOS

Environment=PATH=/home/pi/src/GradientOS/.venv/bin:/home/pi/.local/bin:/home/pi/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/pi/src/GradientOS/src
Environment=GRADIENT_API_HOST=0.0.0.0

ExecStart=/home/pi/src/GradientOS/.venv/bin/gradient-api

Restart=always
RestartSec=10

TimeoutStopSec=30

ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/pi/src/GradientOS /home/pi/.cache /home/pi/.config
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
