# systemd/ethercat-host/ethercat.service.d/10-gradient.conf
[Unit]
# Ensure host tuning runs before EtherCAT master modules load.
Wants=gradient-ethercat-nic-tune.service gradient-irq-affinity.service gradient-cpu-performance.service
After=gradient-ethercat-nic-tune.service gradient-irq-affinity.service gradient-cpu-performance.service

[Service]
# Force IgH ethercatctl to use our pinned config location.
# (IgH's default is /usr/local/etc/ethercat.conf when installed from source.)
ExecStart=
ExecStart=/usr/local/sbin/ethercatctl -c /etc/ethercat.conf start
ExecStop=
ExecStop=/usr/local/sbin/ethercatctl -c /etc/ethercat.conf stop

